# docker/s6/web/run
#!/bin/sh

cd /var/app
gunicorn -w 10 -b 127.0.0.1:8081 pokebox.wsgi
