# docker/s6/notify/run
#!/bin/sh

cd /var/app
python manage.py notify