# docker/crontab
*/5 * * * * python /var/app/manage.py scan
